# Mackup/Library/LaunchAgents/com.amazon.codewhisperer.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.amazon.codewhisperer.launcher</string>
        <key>ProgramArguments</key>
        <array>
            <string>/Applications/Amazon Q.app/Contents/MacOS/q_desktop</string>
            <string>--is-startup</string>
            <string>--no-dashboard</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>AssociatedBundleIdentifiers</key>
        <array>
            <string>com.amazon.codewhisperer</string>
        </array>
    </dict>
</plist>
